# roles/jira/files/mpont_glusterfs.service
Description=jiralb
#After=network.target
#After=docker.service

[Service]
ExecStart=/opt/mount_glfs.sh 

[Install]
WantedBy=multi-user.target
