# roles/keepalived_jiralb/files/keepalived.timer
[Timer]
OnBootSec=30
Unit=keepalived_restart.service
[Install]
WantedBy=multi-user.target
