# roles/keepalived_stor/files/keepalived_restart.service
[Unit]
After=syslog.target network-online.target

[Service]
ExecStart=/opt/keepalived_restart.sh


[Install]
WantedBy=multi-user.target