# roles/pgbouncer/files/pgbouncer.timer
[Timer]
OnBootSec=1min
Unit=pgbouncer.service
[Install]
WantedBy=multi-user.target
