# roles/pgbouncer/files/pgbouncer_start.service
[Unit]
Description=pgbouncer start


[Service]
ExecStart=/opt/pgbouncer_start.sh

[Install]
WantedBy=multi-user.target
