# roles/journald_cli/files/jrnl_u_start.service
[Unit]
Description=Journal upload start


[Service]
ExecStart=/opt/jrnl_u_start.sh

[Install]
WantedBy=multi-user.target
