# roles/journald_cli/files/jrnl_upload.timer
[Timer]
OnBootSec=30
Unit=jrnl_u_start.service
[Install]
WantedBy=multi-user.target
